# systemd/user/pi-codex-session-ping-06-59.timer
[Unit]
Description=Pi Codex session cadence ping at 06:59

[Timer]
OnCalendar=*-*-* 06:59:00
Persistent=true
RandomizedDelaySec=0
AccuracySec=1min
Unit=pi-codex-session-ping.service

[Install]
WantedBy=timers.target

