# systemd/user/pi-codex-session-ping-12-00.timer
[Unit]
Description=Pi Codex session cadence ping at 12:00

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true
RandomizedDelaySec=0
AccuracySec=1min
Unit=pi-codex-session-ping.service

[Install]
WantedBy=timers.target

