# launchd/com.pi-codex-session-ping.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pi-codex-session-ping</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USER/.local/bin/pi-codex-session-ping</string>
  </array>

  <key>StartInterval</key>
  <integer>18000</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/Library/Logs/pi-codex-session-ping.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/Library/Logs/pi-codex-session-ping.err.log</string>
</dict>
</plist>

